# filesystem/opt/ros/melodic/etc/catkin/profile.d/1.ros_os_override.sh
# generated from ros_environment/env-hooks/1.ros_os_override.sh.in

export ROS_OS_OVERRIDE="openembedded"
