# filesystem/var/lib/opkg/info/busybox-syslog.list
/lib/systemd/system/busybox-syslog.service
/etc/syslog.conf
/etc/default/busybox-syslog
/lib/systemd/system/busybox-klogd.service
