# filesystem/var/lib/opkg/info/run-postinsts.list
/lib/
/usr/
/lib/systemd/
/lib/systemd/system/run-postinsts.service
/etc/
/lib/systemd/system/
/
/usr/sbin/
/usr/sbin/run-postinsts
