# filesystem/var/lib/opkg/info/systemd.list
/usr/lib/systemd/user/basic.target
/lib/systemd/system/systemd-fsck@.service
/lib/systemd/system/systemd-reboot.service
/lib/systemd/system/kmod-static-nodes.service
/lib/systemd/system/machines.target.wants/var-lib-machines.mount
/usr/share/dbus-1/system.d/org.freedesktop.timesync1.conf
/lib/systemd/system/timers.target.wants/systemd-tmpfiles-clean.timer
/lib/systemd/system/systemd-exit.service
/usr/share/polkit-1/actions/org.freedesktop.systemd1.policy
/usr/lib/systemd/user/graphical-session.target
/lib/systemd/system/local-fs.target.wants/systemd-remount-fs.service
/usr/bin/systemd-resolve
/lib/systemd/systemd-update-utmp
/lib/systemd/system/runlevel1.target.wants/
/usr/lib/tmpfiles.d/var.conf
/usr/share/polkit-1/actions/org.freedesktop.login1.policy
/usr/bin/resolvectl
/lib/systemd/system/systemd-timesyncd.service
/lib/systemd/system/sysinit.target.wants/sys-fs-fuse-connections.mount
/usr/share/dbus-1/system-services/org.freedesktop.locale1.service
/lib/systemd/system/initrd-parse-etc.service
/lib/systemd/system/dev-mqueue.mount
/lib/systemd/system/initrd.target
/usr/lib/tmpfiles.d/home.conf
/lib/systemd/system/sys-kernel-debug.mount
/usr/lib/systemd/user-preset/90-systemd.preset
/lib/systemd/system/systemd-hibernate-resume@.service
/lib/systemd/systemd-hostnamed
/lib/systemd/systemd-modules-load
/etc/xdg/
/lib/systemd/system/bluetooth.target
/lib/systemd/system/systemd-time-wait-sync.service
/lib/systemd/system/systemd-firstboot.service
/usr/lib/systemd/catalog/systemd.catalog
/sbin/halt
/bin/journalctl
/sbin/reboot
/lib/systemd/system/runlevel4.target
/lib/systemd/system/initrd-switch-root.service
/etc/systemd/system/multi-user.target.wants/
/lib/systemd/system/sysinit.target.wants/systemd-hwdb-update.service
/usr/lib/systemd/user/printer.target
/usr/lib/systemd/catalog/systemd.be@latin.catalog
/usr/share/polkit-1/rules.d/systemd-networkd.rules
/lib/systemd/system/dbus-org.freedesktop.login1.service
/usr/lib/systemd/user/systemd-tmpfiles-clean.service
/lib/systemd/system/initrd-root-device.target
/etc/systemd/system/getty.target.wants/getty@tty1.service
/bin/loginctl
/lib/systemd/systemd-update-done
/etc/systemd/system/sysinit.target.wants/systemd-timesyncd.service
/lib/systemd/system/network.target
/lib/systemd/system/local-fs.target.wants/tmp.mount
/lib/systemd/system/initrd-fs.target
/usr/share/polkit-1/rules.d/
/lib/systemd/system/sysinit.target.wants/systemd-journald.service
/lib/systemd/system/systemd-localed.service
/lib/systemd/system/remote-fs.target.wants/
/lib/systemd/system/sysinit.target.wants/systemd-tmpfiles-setup-dev.service
/usr/share/dbus-1/system-services/org.freedesktop.timesync1.service
/lib/systemd/system/initrd-root-fs.target
/lib/systemd/system/smartcard.target
/lib/systemd/system/ctrl-alt-del.target
/lib/systemd/system/container-getty@.service
/lib/systemd/system/remote-fs.target
/lib/systemd/systemd-sulogin-shell
/usr/lib/systemd/user/graphical-session-pre.target
/lib/systemd/system/sysinit.target.wants/systemd-tmpfiles-setup.service
/usr/share/factory/etc/pam.d/
/lib/systemd/system/poweroff.target.wants/systemd-update-utmp-runlevel.service
/bin/systemd-tmpfiles
/lib/systemd/systemd-timesyncd
/usr/share/dbus-1/system-services/org.freedesktop.systemd1.service
/lib/systemd/system/graphical.target.wants/
/lib/systemd/system/graphical.target.wants/systemd-update-utmp-runlevel.service
/lib/systemd/system/local-fs.target.wants/
/usr/lib/tmpfiles.d/x11.conf
/lib/systemd/system/runlevel0.target
/lib/systemd/system/timers.target.wants/
/lib/systemd/system/systemd-journald.service
/lib/systemd/systemd
/lib/systemd/system/rpcbind.target
/usr/lib/systemd/catalog/systemd.fr.catalog
/lib/systemd/systemd-journald
/lib/systemd/system/rescue.target.wants/
/lib/systemd/system/hibernate.target
/lib/systemd/system/sysinit.target.wants/systemd-sysctl.service
/lib/systemd/system/sysinit.target.wants/systemd-machine-id-commit.service
/usr/share/dbus-1/system.d/org.freedesktop.login1.conf
/lib/systemd/system/systemd-machine-id-commit.service
/lib/systemd/system/sysinit.target.wants/systemd-update-utmp.service
/usr/lib/tmpfiles.d/journal-nocow.conf
/lib/systemd/network/
/lib/systemd/system/shutdown.target
/etc/systemd/system/sysinit.target.wants/
/lib/systemd/system/suspend.target
/etc/pam.d/
/lib/systemd/system/dbus-org.freedesktop.timedate1.service
/lib/systemd/system/sound.target
/lib/systemd/system/user@.service
/usr/share/dbus-1/services/org.freedesktop.systemd1.service
/usr/lib/systemd/user/sound.target
/usr/share/dbus-1/system.d/org.freedesktop.systemd1.conf
/lib/systemd/system/system-update-cleanup.service
/lib/systemd/systemd-backlight
/usr/lib/sysctl.d/50-default.conf
/etc/systemd/system/dbus-org.freedesktop.resolve1.service
/lib/systemd/systemd-dissect
/lib/systemd/system/suspend-then-hibernate.target
/lib/systemd/system/sys-fs-fuse-connections.mount
/lib/systemd/system/runlevel5.target
/etc/systemd/system/dbus-org.freedesktop.network1.service
/lib/systemd/system/sysinit.target.wants/systemd-journal-catalog-update.service
/lib/systemd/system/systemd-modules-load.service
/lib/systemd/system/sockets.target.wants/systemd-journald-dev-log.socket
/lib/systemd/system/ldconfig.service
/lib/systemd/systemd-resolved
/usr/lib/systemd/catalog/systemd.bg.catalog
/usr/bin/busctl
/usr/bin/localectl
/usr/lib/systemd/user/exit.target
/etc/systemd/system/multi-user.target.wants/machines.target
/lib/systemd/system/systemd-tmpfiles-clean.timer
/lib/systemd/system/sysinit.target.wants/systemd-update-done.service
/lib/systemd/systemd-time-wait-sync
/lib/systemd/system/emergency.service
/etc/tmpfiles.d/00-create-volatile.conf
/lib/systemd/system/runlevel4.target.wants/
/usr/lib/sysctl.d/
/usr/lib/systemd/catalog/systemd.ru.catalog
/lib/systemd/systemd-fsck
/lib/systemd/libsystemd-shared-239.so
/lib/systemd/system/systemd-suspend-then-hibernate.service
/lib/systemd/system/paths.target
/etc/systemd/system/network-online.target.wants/systemd-networkd-wait-online.service
/lib/systemd/system/slices.target
/usr/lib/systemd/user/paths.target
/lib/systemd/system/systemd-networkd.socket
/lib/systemd/network/80-container-ve.network
/etc/modules-load.d/
/lib/systemd/system/user-.slice.d/10-defaults.conf
/lib/systemd/system/quotaon.service
/lib/systemd/system/rescue.service
/lib/systemd/system/sysinit.target
/usr/lib/systemd/catalog/systemd.de.catalog
/lib/systemd/system/dbus-org.freedesktop.hostname1.service
/lib/systemd/systemd-quotacheck
/lib/systemd/system/console-getty.service
/lib/systemd/systemd-hibernate-resume
/usr/share/polkit-1/actions/org.freedesktop.timedate1.policy
/lib/systemd/system/sysinit.target.wants/systemd-random-seed.service
/usr/lib/environment.d/99-environment.conf
/usr/lib/systemd/user-preset/
/lib/systemd/system/systemd-journald-dev-log.socket
/lib/systemd/system/systemd-tmpfiles-clean.service
/lib/systemd/system/getty@.service
/lib/systemd/system/exit.target
/lib/systemd/network/80-container-host0.network
/usr/lib/systemd/catalog/systemd.zh_CN.catalog
/lib/systemd/system/runlevel1.target
/lib/systemd/system/systemd-sysctl.service
/lib/systemd/system/multi-user.target.wants/systemd-update-utmp-runlevel.service
/lib/systemd/system/user-runtime-dir@.service
/lib/systemd/system/initrd-cleanup.service
/lib/systemd/systemd-shutdown
/lib/systemd/system/final.target
/lib/systemd/system-generators/systemd-fstab-generator
/lib/systemd/system-generators/systemd-hibernate-resume-generator
/usr/lib/tmpfiles.d/systemd-nologin.conf
/lib/systemd/systemd-networkd-wait-online
/lib/systemd/network/99-default.link
/lib/systemd/system/local-fs.target
/lib/systemd/system/systemd-volatile-root.service
/etc/systemd/resolved.conf
/usr/share/dbus-1/system.d/org.freedesktop.timedate1.conf
/usr/lib/systemd/user/bluetooth.target
/lib/systemd/systemd-timedated
/etc/resolv-conf.systemd
/lib/systemd/system/runlevel2.target.wants/
/lib/systemd/system/systemd-hwdb-update.service
/usr/lib/systemd/catalog/systemd.zh_TW.catalog
/lib/systemd/system/systemd-journald.socket
/usr/lib/systemd/user-environment-generators/
/lib/systemd/system/sysinit.target.wants/systemd-modules-load.service
/usr/share/dbus-1/system.d/org.freedesktop.locale1.conf
/lib/systemd/system/sysinit.target.wants/dev-mqueue.mount
/lib/systemd/systemd-remount-fs
/lib/systemd/system/remote-fs.target.wants/var-lib-machines.mount
/lib/systemd/system/sysinit.target.wants/kmod-static-nodes.service
/lib/systemd/system/sysinit.target.wants/systemd-sysusers.service
/usr/share/polkit-1/
/lib/systemd/system/default.target
/bin/systemd-sysusers
/lib/systemd/system/user.slice
/usr/lib/tmpfiles.d/systemd.conf
/lib/systemd/system/reboot.target
/lib/systemd/system/systemd-journal-catalog-update.service
/usr/share/factory/etc/pam.d/system-auth
/usr/lib/tmpfiles.d/etc.conf
/lib/systemd/systemd-localed
/usr/lib/systemd/catalog/systemd.pt_BR.catalog
/lib/systemd/system/runlevel6.target
/usr/lib/systemd/user/smartcard.target
/lib/systemd/resolv.conf
/lib/systemd/system/dbus-org.freedesktop.locale1.service
/usr/lib/systemd/catalog/systemd.pl.catalog
/lib/systemd/system/systemd-poweroff.service
/usr/share/dbus-1/system.d/org.freedesktop.resolve1.conf
/lib/systemd/system/machines.target.wants/
/lib/systemd/system/poweroff.target.wants/
/usr/lib/tmpfiles.d/tmp.conf
/usr/bin/systemd-mount
/lib/systemd/system-shutdown/
/usr/lib/systemd/user-environment-generators/30-systemd-environment-d-generator
/etc/pam.d/systemd-user
/lib/systemd/system/sys-kernel-config.mount
/sbin/poweroff
/usr/lib/systemd/user/shutdown.target
/lib/systemd/system-generators/systemd-debug-generator
/lib/systemd/system/systemd-resolved.service
/usr/lib/systemd/user/systemd-tmpfiles-setup.service
/etc/sysctl.d/
/usr/lib/systemd/catalog/systemd.it.catalog
/lib/systemd/system/sysinit.target.wants/systemd-firstboot.service
/lib/systemd/system/systemd-networkd.service
/lib/systemd/system/systemd-sysusers.service
/lib/systemd/system/systemd-timedated.service
/usr/share/dbus-1/system-services/org.freedesktop.resolve1.service
/usr/share/systemd/kbd-model-map
/lib/systemd/system/sysinit.target.wants/systemd-journal-flush.service
/lib/systemd/system/tmp.mount
/usr/share/dbus-1/system-services/org.freedesktop.network1.service
/lib/systemd/systemd-user-runtime-dir
/etc/systemd/system/sockets.target.wants/
/etc/systemd/system/multi-user.target.wants/remote-fs.target
/lib/systemd/system/systemd-tmpfiles-setup-dev.service
/usr/share/dbus-1/system.d/org.freedesktop.network1.conf
/usr/share/factory/etc/pam.d/other
/lib/systemd/systemd-growfs
/lib/systemd/system/nss-user-lookup.target
/lib/systemd/system/halt.target
/lib/systemd/system/systemd-backlight@.service
/lib/systemd/systemd-logind
/usr/bin/hostnamectl
/bin/networkctl
/var/lib/systemd/
/usr/lib/sysusers.d/basic.conf
/lib/systemd/system/dev-hugepages.mount
/lib/systemd/system/sockets.target.wants/systemd-journald.socket
/lib/systemd/system/systemd-update-utmp.service
/usr/share/factory/etc/
/lib/systemd/system/swap.target
/lib/systemd/system/systemd-journald-audit.socket
/lib/systemd/network/80-container-vz.network
/lib/systemd/system/network-pre.target
/lib/systemd/system/runlevel2.target
/etc/systemd/timesyncd.conf
/lib/systemd/system/rescue.target
/sbin/telinit
/lib/systemd/system/sigpwr.target
/usr/share/systemd/
/usr/share/dbus-1/system-services/org.freedesktop.timedate1.service
/lib/systemd/system/systemd-networkd-wait-online.service
/lib/systemd/system/systemd-hostnamed.service
/usr/share/polkit-1/actions/org.freedesktop.locale1.policy
/lib/systemd/system/remote-fs-pre.target
/usr/bin/systemd-socket-activate
/usr/lib/systemd/user/default.target
/usr/share/systemd/language-fallback-map
/lib/systemd/system/systemd-halt.service
/usr/share/dbus-1/system-services/org.freedesktop.login1.service
/sbin/shutdown
/usr/share/polkit-1/actions/org.freedesktop.resolve1.policy
/lib/systemd/system/system-update-pre.target
/usr/share/factory/
/etc/xdg/systemd/user
/bin/systemd-notify
/lib/systemd/system/systemd-kexec.service
/lib/systemd/system/rescue.target.wants/systemd-update-utmp-runlevel.service
/lib/systemd/system/systemd-update-done.service
/lib/systemd/system/sockets.target
/etc/systemd/system/multi-user.target.wants/systemd-resolved.service
/lib/systemd/system/getty-pre.target
/lib/systemd/system/sysinit.target.wants/ldconfig.service
/lib/systemd/system/sysinit.target.wants/dev-hugepages.mount
/lib/systemd/system/systemd-rfkill.service
/lib/systemd/system/poweroff.target
/etc/systemd/system/multi-user.target.wants/systemd-networkd.service
/etc/systemd/user/
/lib/systemd/system/graphical.target
/lib/systemd/system/kexec.target
/lib/systemd/systemd-rfkill
/lib/systemd/system/multi-user.target.wants/systemd-logind.service
/lib/systemd/system/systemd-fsck-root.service
/lib/systemd/system/initrd-switch-root.target
/usr/lib/systemd/user/systemd-exit.service
/sbin/runlevel
/usr/bin/systemd-umount
/lib/systemd/system/runlevel5.target.wants/
/lib/systemd/system/systemd-random-seed.service
/usr/share/dbus-1/system-services/org.freedesktop.hostname1.service
/lib/systemd/system/timers.target
/var/volatile/log/journal/
/etc/systemd/system/sockets.target.wants/systemd-networkd.socket
/usr/lib/systemd/catalog/systemd.be.catalog
/etc/systemd/network/
/lib/modprobe.d/systemd.conf
/lib/systemd/system/system-update.target
/lib/systemd/system/user-.slice.d/
/lib/systemd/systemd-networkd
/lib/systemd/systemd-volatile-root
/usr/lib/environment.d/
/lib/systemd/system/syslog.socket
/lib/systemd/system/local-fs-pre.target
/usr/lib/systemd/user/systemd-tmpfiles-clean.timer
/lib/systemd/system/hybrid-sleep.target
/lib/systemd/system/systemd-journal-flush.service
/etc/systemd/system/network-online.target.wants/
/lib/systemd/system/getty.target
/etc/xdg/systemd/
/lib/systemd/system/debug-shell.service
/lib/systemd/system/sockets.target.wants/systemd-journald-audit.socket
/usr/lib/tmpfiles.d/legacy.conf
/lib/systemd/system/multi-user.target
/lib/systemd/system/systemd-rfkill.socket
/lib/systemd/system/systemd-logind.service
/lib/systemd/system/emergency.target
/lib/systemd/system/sysinit.target.wants/sys-kernel-debug.mount
/lib/systemd/system/runlevel3.target
/lib/systemd/system/systemd-tmpfiles-setup.service
/lib/systemd/system/multi-user.target.wants/getty.target
/lib/systemd/system/systemd-quotacheck.service
/lib/systemd/system-generators/systemd-rc-local-generator
/usr/share/dbus-1/system.d/org.freedesktop.hostname1.conf
/lib/systemd/system-preset/
/usr/lib/systemd/user-generators/
/usr/lib/systemd/user/timers.target
/usr/bin/timedatectl
/lib/systemd/system/nss-lookup.target
/lib/systemd/system/network-online.target
/usr/share/factory/etc/nsswitch.conf
/lib/systemd/system/printer.target
/lib/systemd/systemd-makefs
/lib/systemd/system/umount.target
/sbin/init
/bin/systemctl
/lib/systemd/system/runlevel3.target.wants/
/lib/systemd/system/sysinit.target.wants/sys-kernel-config.mount
/usr/lib/modules-load.d/
/lib/systemd/system/systemd-update-utmp-runlevel.service
/lib/systemd/system-generators/systemd-system-update-generator
/sbin/resolvconf
/usr/lib/sysusers.d/systemd.conf
/lib/systemd/systemd-random-seed
/usr/lib/systemd/user/sockets.target
/bin/systemd-machine-id-setup
/usr/share/polkit-1/actions/
/lib/systemd/systemd-sysv-install
/lib/systemd/system-generators/systemd-getty-generator
/usr/lib/systemd/catalog/
/lib/systemd/system/halt-local.service
/lib/systemd/system-preset/90-systemd.preset
/lib/systemd/systemd-sysctl
/usr/share/polkit-1/actions/org.freedesktop.hostname1.policy
/bin/systemd-firstboot
/lib/systemd/system/reboot.target.wants/
/lib/systemd/system/time-sync.target
/lib/systemd/system/reboot.target.wants/systemd-update-utmp-runlevel.service
/lib/systemd/system/autovt@.service
/lib/systemd/system/systemd-remount-fs.service
/lib/systemd/system-generators/systemd-sysv-generator
/lib/systemd/system-generators/
/lib/systemd/system/rc-local.service
/bin/systemd-hwdb
/lib/systemd/system/basic.target
