# filesystem/var/lib/opkg/info/volatile-binds.list
/etc/tmpfiles.d/
/lib/systemd/system/var-volatile-spool.service
/lib/systemd/system/var-volatile-srv.service
/sbin/mount-copybind
/var/cache/
/lib/systemd/system/var-volatile-cache.service
/etc/tmpfiles.d/etc.conf
/lib/systemd/system/var-volatile-lib.service
/etc/tmpfiles.d/home.conf
